# rootfs/etc/services.d/tg-consumer/run
#!/usr/bin/with-contenv sh

python3 /solar-notifier-registrator/run.py